# tests/data/launch_agents/voice_over/com.apple.VoiceOver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.VoiceOver</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.VoiceOver.startup</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
		<key>com.apple.VoiceOver.running</key>
		<dict>
			<key>HideUntilCheckIn</key>
			<true/>
			<key>ResetAtClose</key>
			<true/>
		</dict>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/CoreServices/VoiceOver.app/Contents/MacOS/VoiceOver</string>
		<string>launchd</string>
		<string>-s</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>KeepAlive</key>
	<false/>
	<key>EnableTransactions</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>1</integer>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>Aqua</string>
		<string>LoginWindow</string>
	</array>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
</dict>
</plist>
